# misc/80-pwm-noroot.rules
SUBSYSTEM=="pwm", ACTION=="add|change", \
        RUN+="/bin/chgrp -R plugdev '/sys%p'", \
        RUN+="/bin/chmod -R g=u '/sys%p'"